# versionDist/tasky.service

[Unit]
Description=Tasky

[Service]
ExecStart=CHANGEME/tasky-frontend/tasky

[Install]
WantedBy=multi-user.target 
